# curator4-6648785734.conf
[program:curator4-6648785734]
command=python3 /home/curator4/curator4-6648785734/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator4/curator4-6648785734/test.err.log
stdout_logfile=/home/curator4/curator4-6648785734/test.out.log